# ip_viber/ip-viber.service
[Unit]
Description=IP Viber - WireGuard IP Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ip-viber
ExecStart=/opt/ip-viber/ip_viber
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 